# packages/bitcore-node/bitcore-node.service
[Unit]
Description=Astracore Node
After=network.target

[Service]
User=astracore
Group=astracore
ExecStart=/usr/bin/node /home/astracore/astracore/packages/astracore-node/server.js

[Install]
WantedBy=multi-user.target
